# systemd/.config/systemd/user/disable_screen.timer
[Unit]
Description=Loud time notification

[Timer]
OnCalendar=*-*-* 23:50:00
OnCalendar=*-*-* 23:55:00

[Install]
WantedBy=timers.target
