# systemd/.config/systemd/user/disable_screen.service
[Unit]
Description=disable screen

[Service]
Type=oneshot
ExecStart=/usr/bin/env xset dpms force suspend
