# systemd/.config/systemd/user/disable_screen.service
[Unit]
Description=disable screen (in a bothersome way)

[Service]
Type=oneshot
ExecStart=/usr/bin/env bother
